# services/systemd/swaync.service.in
[Unit]
Description=Swaync notification daemon
Documentation=https://github.com/ErikReider/SwayNotificationCenter
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target
# ConditionEnvironment requires systemd v247 to work correctly
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=dbus
BusName=org.freedesktop.Notifications
ExecStart=@bindir@/swaync
ExecReload=@bindir@/swaync-client --reload-config && @bindir@/swaync-client --reload-css
Restart=on-failure

[Install]
WantedBy=graphical-session.target
